# provisioning/files/mailhog/mailhog.conf
description "MailHog Email Catcher"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /usr/bin/mailhog
